# vmagent/conf/wdstack_agent_v2.conf
[program:wdstack_agent_v2]
command=python /home/op/wdstack_agent_v2/agent.py
autorstart=true
directory=/home/op/wdstack_agent_v2/
autorestart=true
startsecs=10
startretries=10
redirect_stderr=false
